# pacman/.config/systemd/user/notify-updates.service
[Unit]
Description=System update checker

[Service]
Type=oneshot
ExecStart=%h/bin/notify-updates
